# templates/default/monit.conf.erb
# This is an upstart script to keep riak running.
# To install disable the old way of doing things:
#
#   /etc/init.d/riak stop && update-rc.d -f riak remove
#
# then put this script here:
#
#   /etc/init/riak.conf
#
# and reload upstart configuration:
#
#   initctl reload-configuration
#
# You can manually start and stop riak like this:
# 
# start riak
# stop riak
#

description "riak service manager"

limit core unlimited unlimited

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

exec <%= node[:riak][:bin_file] %> -c <%= node[:riak][:riakrc_file] %>

pre-stop exec <%= node[:riak][:bin_file] %> -c <%= node[:riak][:riakrc_file] %> quit